# bt-proxy.service
[Unit]
Description=Bluetooth Proxy (ESPHome-compatible)
After=bluetooth.target dbus.service
Wants=bluetooth.target

[Service]
Type=simple
WorkingDirectory=/opt/bt-proxy
ExecStart=uv run python -m bt_proxy --log-level INFO
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
